# babushka-deps/ghost/init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          ghost
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: ghost initscript
# Description:       ghost
### END INIT INFO

# Original author: Vikhyat Korrapati
 
# Do NOT "set -e"

ROOT=<%= root %>
DAEMON=/usr/local/bin/forever
SCRIPT_NAME=/etc/init.d/ghost
NODE_ENV=production

case "$1" in
  start)
        NODE_ENV=production $DAEMON --append --uid "ghost" --workingDir $ROOT --sourceDir $ROOT start index.js ;;
  stop)
        $DAEMON stop ghost ;;
  restart)
        $SCRIPT_NAME stop && $SCRIPT_NAME start ;;
  status)
        $DAEMON list | grep ghost ;;
  *)
        echo "Usage: $SCRIPT_NAME {start|stop|restart|status}" >&2
        exit 3 ;;
esac

:
